# config/systemd-user/vdirsyncer.service
[Unit]
Description=Calendar synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
