# config/systemd-user/timers.target.wants/vdirsyncer.timer
[Unit]
Description=Calendar synchronization timer

[Timer]
OnBootSec=0m
OnUnitActiveSec=5m
Unit=vdirsyncer.service

[Install]
WantedBy=timers.target
